# llm-proxy.user.service
[Unit]
Description=LLM Tool Call Transform Proxy
After=network.target

[Service]
Type=exec

# Environment configuration
Environment=LLM_PROXY_BACKEND=http://localhost:8080
Environment=LLM_PROXY_PORT=4000
Environment=LLM_PROXY_HOST=127.0.0.1

# Working directory - adjust to where you cloned the repo
WorkingDirectory=%h/llm-retry-proxy

# Run command - choose ONE:

# Option A: With nix-direnv (recommended if using flakes)
ExecStart=/usr/bin/env -C %h/llm-retry-proxy bash -c 'eval "$(direnv export bash)" && python transform_proxy.py --backend ${LLM_PROXY_BACKEND} --port ${LLM_PROXY_PORT} --host ${LLM_PROXY_HOST}'

# Option B: Direct nix run
# ExecStart=%h/.nix-profile/bin/nix run %h/llm-retry-proxy#transform -- --backend ${LLM_PROXY_BACKEND} --port ${LLM_PROXY_PORT} --host ${LLM_PROXY_HOST}

# Option C: Direct Python (if deps installed globally or in venv)
# ExecStart=/usr/bin/python3 %h/llm-retry-proxy/transform_proxy.py --backend ${LLM_PROXY_BACKEND} --port ${LLM_PROXY_PORT} --host ${LLM_PROXY_HOST}

Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
