# toolbridge.service
[Unit]
Description=LLM Tool Call Transform Proxy
Documentation=https://github.com/your-repo/llm-retry-proxy
After=network.target
# If llama-server is also a systemd service, uncomment:
# After=network.target llama-server.service
# Wants=llama-server.service

[Service]
Type=exec
User=nobody
Group=nogroup

# Environment configuration
Environment=TOOLBRIDGE_BACKEND=http://localhost:8080
Environment=TOOLBRIDGE_PORT=4000
Environment=TOOLBRIDGE_HOST=127.0.0.1

# Optional: load from environment file
# EnvironmentFile=/etc/toolbridge/config

# Working directory (adjust to your installation path)
WorkingDirectory=/opt/toolbridge

# Run command - choose ONE of these options:

# Option A: Direct Python (if installed via pip)
ExecStart=/usr/bin/python3 transform_proxy.py --backend ${TOOLBRIDGE_BACKEND} --port ${TOOLBRIDGE_PORT} --host ${TOOLBRIDGE_HOST}

# Option B: Using nix run (if using flake)
# ExecStart=/usr/bin/env nix run .#transform -- --backend ${TOOLBRIDGE_BACKEND} --port ${TOOLBRIDGE_PORT} --host ${TOOLBRIDGE_HOST}

# Option C: Multi-worker production (gunicorn)
# ExecStart=/usr/bin/gunicorn transform_proxy:app -w 4 -k uvicorn.workers.UvicornWorker -b ${TOOLBRIDGE_HOST}:${TOOLBRIDGE_PORT}

# Restart policy
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true

# Resource limits (adjust as needed)
MemoryMax=512M
CPUQuota=100%

[Install]
WantedBy=multi-user.target
